# launchd/com.emailtriage.outlook-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.emailtriage.outlook-monitor</string>
    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/bin/python3</string>
      <string>__REPO_ROOT__/browser/outlook_live_monitor.py</string>
      <string>watch</string>
      <string>--include-pinned</string>
      <string>--interval</string>
      <string>30</string>
      <string>--wake-probe-interval</string>
      <string>2</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>__REPO_ROOT__</string>
    <key>StandardOutPath</key>
    <string>__REPO_ROOT__/shared/outlook_monitor.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO_ROOT__/shared/outlook_monitor.stderr.log</string>
  </dict>
</plist>
